# files/system/usr/etc/skel/.config/autostart/shadowos-setup.desktop
[Desktop Entry]
Type=Application
Name=ShadowOS First Run Setup
Exec=bash /home/$USER/.shadowos-setup.sh
Terminal=false
NoDisplay=true
